# etc_berenger/systemd/system/snap-gtk\x2dcommon\x2dthemes-1440.mount
[Unit]
Description=Mount unit for gtk-common-themes, revision 1440
Before=snapd.service

[Mount]
What=/var/lib/snapd/snaps/gtk-common-themes_1440.snap
Where=/snap/gtk-common-themes/1440
Type=squashfs
Options=nodev,ro,x-gdu.hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
